# rod_c/systemd/rod-communication.service
[Unit]
Description=ROD Communication - IPC Thread
After=network.target rod-detection.service
Requires=rod-detection.service
BindsTo=rod-detection.service

[Service]
Type=simple
User=noegame
Group=noegame
WorkingDirectory=/home/noegame/ROD

# Main communication process
ExecStart=/home/noegame/ROD/rod_c/build/rod_communication

# Wait for detection socket to be ready
ExecStartPre=/bin/sleep 1

# Restart policy
Restart=always
RestartSec=2

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rod-communication

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
