# rod_c/systemd/rod.target
[Unit]
Description=ROD - Remote Observation Device (Complete System)
After=network.target

# This target groups both services together
Wants=rod-detection.service rod-communication.service

[Install]
WantedBy=multi-user.target
